# etc/systemd/system/app.service
[Unit]
Description=app

[Service]
Type=simple
Restart=always
ExecStart=/ProteiTestCase

[Install]
WantedBy=multi-user.target